# deploy/smart_agent.service
# ============================================
# SmartAgent systemd 服务文件
# ============================================
#
# 部署步骤:
#   1. sudo cp deploy/smart_agent.service /etc/systemd/system/
#   2. sudo systemctl daemon-reload
#   3. sudo systemctl enable smart_agent
#   4. sudo systemctl start smart_agent
#   5. sudo systemctl status smart_agent
#
# 常用命令:
#   sudo systemctl restart smart_agent    # 重启
#   sudo systemctl stop smart_agent       # 停止
#   sudo journalctl -u smart_agent -f     # 查看日志

[Unit]
Description=SmartAgent - AI Agent Service (Enterprise v2.0)
Documentation=https://github.com/TYYHeck/smart_agent
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/smart_agent

# 环境变量
EnvironmentFile=/root/smart_agent/.env
Environment="PATH=/root/smart_agent/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONWARNINGS=ignore::DeprecationWarning"

# 启动命令
ExecStart=/root/smart_agent/venv/bin/python main.py --web --host 0.0.0.0 --port 8080

# 重启策略
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# 进程管理
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# 日志 (journald)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=smart_agent

# 资源限制
LimitNOFILE=65536
MemoryHigh=2G
MemoryMax=4G
CPUQuota=200%

# 注：ProtectHome/ProtectSystem 等 namespace 隔离与 venv 符号链接不兼容，
# 安全通过 nginx 反向代理 + 防火墙保障

[Install]
WantedBy=multi-user.target
